# resources/.runit/sv/psd/finish
#!/bin/sh

exec psd unsync ; rm -rf ./supervise
